# nats/supervisord.conf
[supervisord]
nodaemon=true

[program:gnatsd]
command=/go/bin/gnatsd &
