# gorouter/supervisord.conf
[supervisord]
nodaemon=true

[program:gorouter]
command=/go/bin/gorouter -c gorouter.yml
